# wsl-vpnkit.service
[Unit]
Description=wsl-vpnkit
After=network.target

[Service]
# Use local wsl-vpnkit binary when available, otherwise fall back to distro invocation.
ExecStart=/bin/sh -c '\
if command -v wsl-vpnkit >/dev/null 2>&1; then \
	exec wsl-vpnkit; \
else \
	exec /mnt/c/Windows/system32/wsl.exe -d wsl-vpnkit --cd /app wsl-vpnkit; \
fi'

Restart=always
KillMode=mixed

[Install]
WantedBy=multi-user.target
